# ops/launchd/com.tennislevel.daily-progress.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tennislevel.daily-progress</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO_ROOT__/scripts/run-daily-progress.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__</string>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>TZ</key>
    <string>America/Edmonton</string>
    <key>CODEX_BIN</key>
    <string>__CODEX_BIN__</string>
    <key>CODEX_LOG_DIR</key>
    <string>__LOG_DIR__</string>
  </dict>
  <key>StandardOutPath</key>
  <string>__LOG_DIR__/launchd-daily.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/launchd-daily.stderr.log</string>
</dict>
</plist>
